# config/services/acme.service
[Unit]
Description=ACME Certificate Renewal (DayShield)
Documentation=man:certbot(1)
After=network-online.target
Wants=network-online.target
ConditionKernelCommandLine=!installer
ConditionPathExists=/etc/letsencrypt

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew --quiet --no-self-upgrade
ExecStartPost=/usr/bin/install -d -m 700 /etc/dayshield/certs
ExecStartPost=/bin/cp -a /etc/letsencrypt/live/. /etc/dayshield/certs/
ExecStartPost=/bin/chmod -R go-rwx /etc/dayshield/certs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
