# config/services/acme.timer
[Unit]
Description=Run DayShield ACME certificate renewal
Documentation=man:systemd.timer(5)
ConditionKernelCommandLine=!installer

[Timer]
OnBootSec=15m
OnUnitActiveSec=24h
Persistent=true
Unit=acme.service

[Install]
WantedBy=timers.target
