# config/services/crowdsec.service
[Unit]
Description=CrowdSec Security Engine
Documentation=https://docs.crowdsec.net
After=network-online.target nftables.service
Wants=network-online.target

[Service]
Type=notify
Environment=LC_ALL=C.UTF-8
EnvironmentFile=-/etc/default/crowdsec
ExecStart=/usr/bin/crowdsec -c /etc/crowdsec/config.yaml
ExecStartPost=/usr/bin/cscli hub update --force 2>/dev/null || true
Restart=on-failure
RestartSec=10s
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/crowdsec /var/log/crowdsec /run/crowdsec /etc/crowdsec

[Install]
WantedBy=multi-user.target
