# config/services/dayshield.service
[Unit]
Description=DayShield Firewall Core
Documentation=https://github.com/daygle/dayshield
After=network-online.target nftables.service unbound.service
Wants=network-online.target
Requires=nftables.service

[Service]
Type=exec
Environment=DAYSHIELD_PORT=8443
ExecStart=/usr/local/sbin/dayshield-core
Restart=on-failure
RestartSec=5s
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/dayshield /var/lib/dayshield /var/log/dayshield /opt/dayshield-core /opt/dayshield-ui /opt/dayshield-rootfs /usr/local/sbin /usr/local/lib/dayshield /usr/local/share
PrivateTmp=yes
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
