# config/services/unbound.service
[Unit]
Description=Unbound DNS Resolver
Documentation=man:unbound(8)
After=network-online.target
Wants=network-online.target
Before=nss-lookup.target
ConditionKernelCommandLine=!installer

[Service]
Type=notify
# Bootstrap DNSSEC trust anchor on first boot (noop if already present)
ExecStartPre=/bin/sh -c 'mkdir -p /var/lib/unbound && /usr/sbin/unbound-anchor -a /var/lib/unbound/root.key || true'
ExecStartPre=/bin/sh -c 'chown -R unbound:unbound /var/lib/unbound 2>/dev/null || true'
ExecStartPre=/usr/sbin/unbound-checkconf
ExecStart=/usr/sbin/unbound -d
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
RuntimeDirectory=unbound
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/unbound /run/unbound
PrivateTmp=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID

[Install]
WantedBy=multi-user.target
